# backend/theta-webui-updater.service
[Unit]
Description=Theta WebUI Auto-Updater (Git Pull)
After=network.target

[Service]
Type=oneshot
User=cloud
WorkingDirectory=/home/cloud/theta-webui/backend

# Check for updates and restart backend if changes detected
ExecStart=/usr/bin/bash -c '\
  cd /home/cloud/theta-webui/backend && \
  git fetch --all  && \
  LOCAL=$(git rev-parse HEAD) && \
  REMOTE=$(git rev-parse @{u}) && \
  if [ "$LOCAL" != "$REMOTE" ]; then \
    echo "Updates detected: $LOCAL -> $REMOTE"; \
    git pull && \
    systemctl restart theta-webui-backend.service && \
    echo "Backend restarted with new changes"; \
  else \
    echo "No updates available"; \
  fi'

StandardOutput=journal
StandardError=journal
SyslogIdentifier=theta-updater
